# supervisord.conf
[supervisord]
loglevel=debug

[program:document-service]
command=node server.js
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
